# root/etc/periodic/daily/get-zone
#!/bin/bash

echo "`date`: this script gets executed daily" >> /tmp/cron.log


